# container-files/etc/supervisor.d/nginx.conf
[program:nginx]
command = /usr/sbin/nginx
autostart = true
autorestart = true
stderr_logfile = NONE
stdout_logfile = NONE
